# motd_banner/update-motd.d/01-message
#!/bin/bash

echo "Message d'informations"
