# motd_banner/update-motd.d/01-message
#!/bin/bash

echo ""
echo "***********************************************"
echo "*                                             *"
echo "*   Access forbidden to unauthorized users.   *"
echo "*                                             *"
echo "***********************************************"
echo ""
