# motd_banner/update-motd.d/02-ip-address
#!/bin/bash

echo "IP addresses of this server are :"

IFS=$'\ '
for address in $(hostname -I)    
do
    echo "$address"
done
